# deploy/systemd/omada-webhook.service
[Unit]
Description=Omada webhook payload capture receiver
After=network.target

[Service]
Type=simple
User=omada
Group=omada
EnvironmentFile=/etc/omada-webhook.env
ExecStart=/usr/bin/python3 /opt/omada-webhook/webhook_server.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/omada-webhook
ReadWritePaths=/var/lib/omada-webhook

[Install]
WantedBy=multi-user.target
